# upstart_unicorn.conf
description "upstart-socket-bridge unicorn"
start on socket PROTO=unix PATH=/var/run/unicorn.sock

exec /usr/bin/ruby /var/www/project/current/unicorn_wrapper.rb
